# packaging/systemd/keydviz-helperd.service.d/keypresses.conf
# Opt-in drop-in: enable keypress glow (live per-key highlight) in addition to layers.
#
# This widens the daemon's access to /dev/input — i.e. it becomes able to read keystrokes
# (still only the foreground desktop user's, still read-only, still sandboxed and
# network-isolated). Install it deliberately:
#
#   sudo install -Dm644 packaging/systemd/keydviz-helperd.service.d/keypresses.conf \
#       /etc/systemd/system/keydviz-helperd.service.d/keypresses.conf
#   sudo systemctl daemon-reload && sudo systemctl restart keydviz-helperd
#
# Remove the file and restart to drop back to layers-only.

[Service]
# With --keys the daemon reads keyd's virtual keyboard from /dev/input/event* directly
# (evdev, in-process — no `keyd monitor` child), so it needs the `input` group and access
# to the input device class. We open the nodes read-only, so `r` is sufficient.
SupplementaryGroups=input
DeviceAllow=char-input r

# Re-issue ExecStart with --keys (a drop-in must clear the inherited one first).
ExecStart=
ExecStart=/usr/bin/keydviz-helperd --active-session --keys --socket /run/keyd-viz/keyd-viz.sock
